# nm-elevation.service
[Unit]
Description=New Mexico Elevation Data Collection Service
After=network.target

[Service]
Type=simple
User=raspberry
WorkingDirectory=/home/raspberry/projects/nm_elevation
ExecStart=/home/raspberry/projects/nm_elevation/collect_service.sh
Restart=always
RestartSec=10
StandardOutput=append:/home/raspberry/projects/nm_elevation/collection.log
StandardError=append:/home/raspberry/projects/nm_elevation/collection.log
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target 